# systemd/docker-maintenance.service
[Unit]
Description=Docker maintenance (rclone backup + container update)
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=berduf
EnvironmentFile=-/home/berduf/.config/docker-maintenance.env
ExecStart=/usr/bin/docker compose \
    -f /home/berduf/docker-idempotent-update/docker-compose.yml \
    run --rm maintenance
